# user-workspace/freeroot/var/lib/dpkg/info/supervisor.conffiles
/etc/default/supervisor
/etc/init.d/supervisor
/etc/supervisor/supervisord.conf
