# user-workspace/freeroot/var/lib/dpkg/info/supervisor.list
/.
/etc
/etc/default
/etc/default/supervisor
/etc/init.d
/etc/init.d/supervisor
/etc/supervisor
/etc/supervisor/conf.d
/etc/supervisor/supervisord.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/supervisor.service
/usr
/usr/bin
/usr/bin/echo_supervisord_conf
/usr/bin/pidproxy
/usr/bin/supervisorctl
/usr/bin/supervisord
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/supervisor
/usr/lib/python3/dist-packages/supervisor/__init__.py
/usr/lib/python3/dist-packages/supervisor/childutils.py
/usr/lib/python3/dist-packages/supervisor/compat.py
/usr/lib/python3/dist-packages/supervisor/confecho.py
/usr/lib/python3/dist-packages/supervisor/datatypes.py
/usr/lib/python3/dist-packages/supervisor/dispatchers.py
/usr/lib/python3/dist-packages/supervisor/events.py
/usr/lib/python3/dist-packages/supervisor/http.py
/usr/lib/python3/dist-packages/supervisor/http_client.py
/usr/lib/python3/dist-packages/supervisor/loggers.py
/usr/lib/python3/dist-packages/supervisor/medusa
/usr/lib/python3/dist-packages/supervisor/medusa/__init__.py
/usr/lib/python3/dist-packages/supervisor/medusa/asynchat_25.py
/usr/lib/python3/dist-packages/supervisor/medusa/asyncore_25.py
/usr/lib/python3/dist-packages/supervisor/medusa/auth_handler.py
/usr/lib/python3/dist-packages/supervisor/medusa/counter.py
/usr/lib/python3/dist-packages/supervisor/medusa/default_handler.py
/usr/lib/python3/dist-packages/supervisor/medusa/filesys.py
/usr/lib/python3/dist-packages/supervisor/medusa/http_date.py
/usr/lib/python3/dist-packages/supervisor/medusa/http_server.py
/usr/lib/python3/dist-packages/supervisor/medusa/logger.py
/usr/lib/python3/dist-packages/supervisor/medusa/producers.py
/usr/lib/python3/dist-packages/supervisor/medusa/util.py
/usr/lib/python3/dist-packages/supervisor/medusa/xmlrpc_handler.py
/usr/lib/python3/dist-packages/supervisor/options.py
/usr/lib/python3/dist-packages/supervisor/pidproxy.py
/usr/lib/python3/dist-packages/supervisor/poller.py
/usr/lib/python3/dist-packages/supervisor/process.py
/usr/lib/python3/dist-packages/supervisor/rpcinterface.py
/usr/lib/python3/dist-packages/supervisor/scripts
/usr/lib/python3/dist-packages/supervisor/scripts/loop_eventgen.py
/usr/lib/python3/dist-packages/supervisor/scripts/loop_listener.py
/usr/lib/python3/dist-packages/supervisor/scripts/sample_commevent.py
/usr/lib/python3/dist-packages/supervisor/scripts/sample_eventlistener.py
/usr/lib/python3/dist-packages/supervisor/scripts/sample_exiting_eventlistener.py
/usr/lib/python3/dist-packages/supervisor/skel
/usr/lib/python3/dist-packages/supervisor/skel/sample.conf
/usr/lib/python3/dist-packages/supervisor/socket_manager.py
/usr/lib/python3/dist-packages/supervisor/states.py
/usr/lib/python3/dist-packages/supervisor/supervisorctl.py
/usr/lib/python3/dist-packages/supervisor/supervisord.py
/usr/lib/python3/dist-packages/supervisor/templating.py
/usr/lib/python3/dist-packages/supervisor/tests
/usr/lib/python3/dist-packages/supervisor/tests/__init__.py
/usr/lib/python3/dist-packages/supervisor/tests/base.py
/usr/lib/python3/dist-packages/supervisor/tests/fixtures
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/donothing.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/example
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/example/included.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/include.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-1054.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-1224.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-1231.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-565.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-638.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-663.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-664.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-835.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/issue-836.conf
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/listener.py
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/spew.py
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/test_1231.py
/usr/lib/python3/dist-packages/supervisor/tests/fixtures/unkillable_spew.py
/usr/lib/python3/dist-packages/supervisor/tests/test_childutils.py
/usr/lib/python3/dist-packages/supervisor/tests/test_confecho.py
/usr/lib/python3/dist-packages/supervisor/tests/test_datatypes.py
/usr/lib/python3/dist-packages/supervisor/tests/test_dispatchers.py
/usr/lib/python3/dist-packages/supervisor/tests/test_end_to_end.py
/usr/lib/python3/dist-packages/supervisor/tests/test_events.py
/usr/lib/python3/dist-packages/supervisor/tests/test_http.py
/usr/lib/python3/dist-packages/supervisor/tests/test_http_client.py
/usr/lib/python3/dist-packages/supervisor/tests/test_loggers.py
/usr/lib/python3/dist-packages/supervisor/tests/test_options.py
/usr/lib/python3/dist-packages/supervisor/tests/test_poller.py
/usr/lib/python3/dist-packages/supervisor/tests/test_process.py
/usr/lib/python3/dist-packages/supervisor/tests/test_rpcinterfaces.py
/usr/lib/python3/dist-packages/supervisor/tests/test_socket_manager.py
/usr/lib/python3/dist-packages/supervisor/tests/test_states.py
/usr/lib/python3/dist-packages/supervisor/tests/test_supervisorctl.py
/usr/lib/python3/dist-packages/supervisor/tests/test_supervisord.py
/usr/lib/python3/dist-packages/supervisor/tests/test_templating.py
/usr/lib/python3/dist-packages/supervisor/tests/test_web.py
/usr/lib/python3/dist-packages/supervisor/tests/test_xmlrpc.py
/usr/lib/python3/dist-packages/supervisor/ui
/usr/lib/python3/dist-packages/supervisor/ui/images
/usr/lib/python3/dist-packages/supervisor/ui/images/button_refresh.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/button_restart.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/button_stop.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/icon.png
/usr/lib/python3/dist-packages/supervisor/ui/images/rule.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/state0.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/state1.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/state2.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/state3.gif
/usr/lib/python3/dist-packages/supervisor/ui/images/supervisor.gif
/usr/lib/python3/dist-packages/supervisor/ui/status.html
/usr/lib/python3/dist-packages/supervisor/ui/stylesheets
/usr/lib/python3/dist-packages/supervisor/ui/stylesheets/supervisor.css
/usr/lib/python3/dist-packages/supervisor/ui/tail.html
/usr/lib/python3/dist-packages/supervisor/version.txt
/usr/lib/python3/dist-packages/supervisor/web.py
/usr/lib/python3/dist-packages/supervisor/xmlrpc.py
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/entry_points.txt
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/not-zip-safe
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/requires.txt
/usr/lib/python3/dist-packages/supervisor-4.1.0.egg-info/top_level.txt
/usr/share
/usr/share/doc
/usr/share/doc/supervisor
/usr/share/doc/supervisor/README.Debian
/usr/share/doc/supervisor/changelog.Debian.gz
/usr/share/doc/supervisor/copyright
/usr/share/doc/supervisor/examples
/usr/share/doc/supervisor/examples/loop_eventgen.py
/usr/share/doc/supervisor/examples/loop_listener.py
/usr/share/doc/supervisor/examples/sample.conf
/usr/share/doc/supervisor/examples/sample_commevent.py
/usr/share/doc/supervisor/examples/sample_eventlistener.py
/usr/share/doc/supervisor/examples/sample_exiting_eventlistener.py
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/echo_supervisord_conf.1.gz
/usr/share/man/man1/pidproxy.1.gz
/usr/share/man/man1/supervisor.1.gz
/usr/share/man/man1/supervisorctl.1.gz
/usr/share/man/man1/supervisord.1.gz
/var
/var/log
/var/log/supervisor
